# deploy/supervisor/all-in-one.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql_err.log
priority=10

[program:redis]
command=/usr/bin/redis-server --requirepass %(ENV_REDIS_PASSWORD)s
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis_err.log
priority=10

[program:qdrant]
command=/usr/local/bin/qdrant --config-path /etc/qdrant/config.yaml
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/qdrant.log
stderr_logfile=/var/log/supervisor/qdrant_err.log
priority=10

[program:backend]
command=python /app/main.py server --host 127.0.0.1 --port 8000 --no-reload
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=10
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_err.log
priority=20

[program:worker]
command=python /app/main.py worker -c 2 -Q default,workflow
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=10
stdout_logfile=/var/log/supervisor/worker.log
stderr_logfile=/var/log/supervisor/worker_err.log
priority=30

[program:beat]
command=python /app/main.py beat
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=10
stdout_logfile=/var/log/supervisor/beat.log
stderr_logfile=/var/log/supervisor/beat_err.log
priority=30

[program:frontend]
command=node /app/frontend/server.js
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend_err.log
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0"
priority=40

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
priority=50
