# deploy/supervisor/app.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:backend]
command=python /app/main.py server --host 0.0.0.0 --port 8000 --no-reload
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_err.log
priority=10

[program:worker]
command=python /app/main.py worker -c 4 -Q default,workflow
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/worker.log
stderr_logfile=/var/log/supervisor/worker_err.log
priority=20

[program:beat]
command=python /app/main.py beat
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/beat.log
stderr_logfile=/var/log/supervisor/beat_err.log
priority=20

[program:frontend]
command=node /app/frontend/server.js
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend_err.log
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0"
priority=30
